# os/stage-camplayer/01-camplayer-install/files/camplayer-firstboot.service
[Unit]
Description=Camplayer first-boot setup
After=network-online.target
Wants=network-online.target
ConditionPathExists=/usr/local/bin/camplayer-firstboot.sh

[Service]
Type=oneshot
ExecStart=/usr/local/bin/camplayer-firstboot.sh
RemainAfterExit=yes
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
